# templates/FreeBSD/cloud-password-rc.j2
#!/bin/sh

. /etc/rc.subr

name=cloudstack-set-guest-password
rcvar=cloudstack_set_guest_password_enable

command="/usr/share/cloudstack/${name}"

load_rc_config $name
run_rc_command "$1"