# templates/firstboot.upstart.j2
# cloudstack-set-guest-sshkey

description "set root sshkeys from CloudStack metadata"

start on startup

task
exec /etc/init.d/firstboot-cloudstack
