# templates/run-userdata.upstart.j2
# cloudstack-set-guest-sshkey

description "run userdata from CloudStack metadata"

start on (startup and
          started networking)

task
exec /usr/share/cloudstack/cloudstack-run-userdata start
