# templates/set-guest-sshkey.upstart.j2
# cloudstack-set-guest-sshkey

description "set root sshkeys from CloudStack metadata"

start on started network-services

script
  /usr/share/cloudstack/cloudstack-set-guest-sshkey start
end script
